# target/linux/mediatek/filogic/base-files/etc/init.d/modem_reboot
#!/bin/sh /etc/rc.common
START=99
STOP=0

start() {
    echo "do nothing"
}

stop() {
    for modem_no in 0 1 2; do
        dev=$(uci -q get modem.modem${modem_no}.at_port)
        [ -n "$dev" ] || continue         
        [ -e "$dev" ] || continue          

        case "$dev" in
            *mhi*)
                echo "[$dev] is an mhi device, skipping"
                continue
                ;;
            *) ;;
        esac

        sms_tool -d $dev at AT+CFUN=1,1
        for i in $(seq 1 15); do
            if [ ! -e "$dev" ]; then
                echo "[$dev] Device disconnected after $i s"
                break
            fi
            sleep 1
        done
    done
}
